# SRC/Hardenwing-OS/config/includes.chroot/etc/systemd/system/usbguard.service.d/30_hw.conf
[Unit]
# Do not fail if the system does not have a USB bus (e.g., in some VMs or minimal servers)
ConditionPathExists=/sys/bus/usb
