# SRC/Overkill/config/includes.chroot/etc/xdg/autostart/conky-reanchor.desktop
[Desktop Entry]
Type=Application
Name=Conky Reanchor
Comment=Auto-reanchor Conky to top-right when resolution/monitor layout changes
Exec=/bin/sh -lc 'sleep 6 && bash /usr/bin/conky-reanchor &'
X-KDE-Autostart-enabled=true
TryExec=/bin/bash
OnlyShowIn=KDE;
NoDisplay=false
Terminal=false