# SRC/ARCHIVED/Slarpx/config/includes.chroot/etc/NetworkManager/dispatcher.d/99-slarpx-dns
#!/bin/sh
# Slarpx DNS Dispatcher
# Ensures DNS remains locked even after NetworkManager events

case "$2" in
    up|down|vpn-up|vpn-down|hostname)
        /usr/bin/slarpx-dns
        ;;
esac
